# ifinance.crontab
SHELL=/bin/bash
MAILTO=xxx@gmail.com
### ifinance
0 19 * * * cd /var/apps/ifinance && bin/ifinance collect >> logs/collect_cron.log 2>&1
0 22 * * * cd /var/apps/ifinance && bin/ifinance report
