# ifinance.crontab
SHELL=/bin/bash
MAILTO=xxx@gmail.com
LANG=ja_JP.UTF-8
CONTENT_TYPE=text/plain; charset="utf-8"
### ifinance
0 19 * * * cd /var/apps/ifinance && bin/ifinance collect >> logs/collect_cron.log 2>&1
0 22 * * * cd /var/apps/ifinance && bin/ifinance report
