# ifinance.crontab
SHELL=/bin/bash
MAILTO=xxx@gmail.com
LANG=ja_JP.UTF-8
CONTENT_TYPE=text/plain; charset="utf-8"
### ifinance
0 16 * * 1-5 cd /var/apps/ifinance && bin/ifinance collect >> logs/collect_cron.log 2>&1
0 23 * * 1-5 cd /var/apps/ifinance && bin/ifinance predict
10 23 * * 1-5 cd /var/apps/ifinance && bin/ifinance report
15 23 * * 1-5 cd /var/apps/ifinance && bin/ifinance report Overvalued
